# systemd/sos-memory.service
[Unit]
Description=SOS Memory Service (Vector DB & Retrieval)
After=network.target

[Service]
Type=simple
User=mumega
Group=mumega
WorkingDirectory=/home/mumega/SOS
Environment=PYTHONPATH=/home/mumega/SOS
Environment=SOS_MEMORY_PORT=6061
EnvironmentFile=-/home/mumega/SOS/.env

ExecStart=/usr/bin/python3 -m sos.services.memory

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sos-memory

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
