# systemd/sos-gateway-bridge.service
[Unit]
Description=SOS Gateway Bridge API (External Agents)
After=network.target sos-engine.service

[Service]
Type=simple
User=mumega
Group=mumega
WorkingDirectory=/home/mumega/SOS
Environment=PYTHONPATH=/home/mumega/SOS
Environment=SOS_BRIDGE_PORT=8001
EnvironmentFile=-/home/mumega/SOS/.env

ExecStart=/usr/bin/python3 -m sos.services.gateway.app

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sos-gateway-bridge

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
