# systemd/sos-engine.service
[Unit]
Description=SOS Engine Service (Orchestration & Multi-Model Routing)
After=network.target

[Service]
Type=simple
User=mumega
Group=mumega
WorkingDirectory=/home/mumega/SOS
Environment=PYTHONPATH=/home/mumega/SOS
Environment=SOS_ENGINE_PORT=6060
# API keys loaded from .env file (not committed)
EnvironmentFile=-/home/mumega/SOS/.env

ExecStart=/usr/bin/python3 -m sos.services.engine

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=sos-engine

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
